# deploy/systemd/autopilot-jobhunt.service
[Unit]
Description=autopilot-jobhunt scheduler service
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/var/www/autopilot-jobhunt
Environment=AUTOPILOT_LOG_FILE=/var/www/autopilot-jobhunt/scan.log
Environment=PYTHONIOENCODING=utf-8
ExecStart=/var/www/autopilot-jobhunt/venv/bin/autopilot service
Restart=always
RestartSec=10
Nice=10

[Install]
WantedBy=multi-user.target
